# tftpd/supervisord-tftpd.conf
[include]
files = base.conf rsyslogd.conf

[program:tftpd]
command=/usr/sbin/in.tftpd --foreground --user tftp --secure /var/spool/tftp
