# httpd/supervisord-httpd.conf
[include]
files = base.conf

[program:httpd]
command=/usr/sbin/apache2ctl -D FOREGROUND
environment=APACHE_LOG_DIR="/var/log/supervisor"
stopasgroup=true
