# supervisord.conf.d/dhcpd.conf
[include]
files = /etc/supervisor/conf.d/base.conf /etc/supervisor/conf.d/rsyslogd.conf

[program:dhcpd]
command=/usr/sbin/dhcpd -q -d -cf /etc/dhcp/dhcpd.conf -pf /var/run/dhcpd.pid eth0
