# supervisord.conf.d/squid.conf
[include]
files = /etc/supervisor/conf.d/base.conf

[program:squid]
command=/usr/sbin/squid3 -YC -f /etc/squid-deb-proxy/squid-deb-proxy.conf -N
#PIDFILE /var/run/squid-deb-proxy.pid
